# philosophie.ch_legacy/sysadmin/crontab
# m h  dom mon dow   command

# To clear the rails cache every day at midnight
0 0 * * * /home/sysadmin/gitrepos/sysadmin-utils/philosophie.ch_legacy/sysadmin/system-tasks/rails_cache_cleanup.sh